# templates/etc/init/mmonit.conf
# M/Monit script for upstart (ubuntu)
#

description "M/Monit daemon"

start on runlevel [2345]
stop on runlevel [06]

exec sudo -u {{ mmonit_user }} mmonit -i
respawn
